# templates/gitlab_runner.erb
#!/bin/sh
#
# This file is managed via Puppet

daemon="/usr/local/bin/gitlab-runner"
daemon_flags="--log-level <%= @log_level %> run"
daemon_user="<%= @user %>"

. /etc/rc.d/rc.subr

pexp="${daemon}${daemon_flags:+ ${daemon_flags}}"
rc_reload=NO

rc_start() {
        ${rcexec} "${daemon} ${daemon_flags}" \
                        2>><%= @log_file %> \
                        >><%= @log_file %> &
}

rc_pre() {
	install -d -o <%= @user %> `dirname <%= @log_file %>`
}

rc_cmd $1
